# exp-monitoring.service.example
[Unit]
Description=exp-monitoring server
Documentation=exp-monitoring - Monitor EXP from Github.
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/bin/python3 /usr/share/EXP-monitoring/EXP-Monitoring.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
Alias=exp-monitoring.service